# templates/default/upstart/eventstore.conf.erb
pre-start script

bash << "EOF"
  mkdir -p <%= @logs_dir %>
  chown -R <%= @user %> <%= @logs_dir %>
  mkdir -p <%= @data_dir %>
  chown -R <%= @user %> <%= @data_dir %>
EOF

end script

start on (started network-interface
          or started network-manager
          or started networking)

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking)

respawn

exec su - <%= @user %> -c 'cd <%= @install_dir %>; <%= @command %> --config=<%= @config_file %>'
