# templates/default/upstart/eventstore.conf.erb
pre-start script

bash << "EOF"
  mkdir -p /var/log/eventstore
  chown -R eventstore /var/log/eventstore
EOF

end script

start on (started network-interface
          or started network-manager
          or started networking)

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking)

respawn

exec su - seif -c 'cd /usr/local/eventstore; mono EventStore.SingleNode.exe >> /var/log/eventstore/eventstore.log 2>&1'
