# templates/default/upstart/eventstore.conf.erb
pre-start script

bash << "EOF"
  mkdir -p <%= node['eventstore']['log_dir'] %>
  chown -R eventstore <%= node['eventstore']['log_dir'] %>
EOF

end script

start on (started network-interface
          or started network-manager
          or started networking)

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking)

respawn

exec su - eventstore -c 'cd <%= node['eventstore']['install_dir'] %>; mono EventStore.SingleNode.exe %> >> <%= node['eventstore']['log_dir'] %>upstart.log 2>&1'
